# ansible/files/nas-01/snapraid-btrfs-runner.timer
[Unit]
Description=Run snapraid-btrfs-runner every night

[Timer]
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
