# ansible/roles/power_saving/templates/power-saving-hard-drives.service.j2
[Unit]
Description=Configure power management for hard drives

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot

ExecStart={{ power_saving_hard_drives_script_path }} configure
