# ansible/roles/syncoid_source/templates/syncoid-dataset.service.j2
[Unit]
Description=syncoid sync from {{ syncoid_source_dataset }} to {{ syncoid_destination_host }}/{{ syncoid_destination_dataset }}
{% if syncoid_healthchecks_uuid is defined %}
OnFailure=healthchecks@%N:{{ syncoid_healthchecks_uuid }}:failure.service
OnSuccess=healthchecks@%N:{{ syncoid_healthchecks_uuid }}:success.service
Wants=healthchecks@%N:{{ syncoid_healthchecks_uuid }}:start.service
{% endif %}
After=network.target

[Service]
User={{ syncoid_source_username }}
Type=oneshot
ExecStart=/usr/sbin/syncoid --no-privilege-elevation --recursive {{ syncoid_source_dataset }} {{ syncoid_destination_username }}@{{ syncoid_destination_host }}:{{ syncoid_destination_dataset }}
